# com.xtm.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xtm.snapshot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jayjay5032/MartinMonthlyReport/generate_report.py</string>
        <string>--snapshot</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jayjay5032/MartinMonthlyReport</string>

    <!-- Run every day at 6:00 PM to capture per-user stats before projects are archived -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jayjay5032/MartinMonthlyReport/xtm_snapshot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jayjay5032/MartinMonthlyReport/xtm_snapshot_error.log</string>

    <!-- Do not run immediately on load -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Set nice value to run with lower priority -->
    <key>Nice</key>
    <integer>1</integer>

    <!-- Set process timeout to 30 minutes -->
    <key>ExitTimeOut</key>
    <integer>1800</integer>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
